# deploy/systemd/pandapages-postgresql-backup-healthcheck.service
[Unit]
Description=Detect stale Panda Pages PostgreSQL backup verification state
OnFailure=pandapages-postgresql-backup-notify@%n.service

[Service]
Type=oneshot
User=root
Group=root
EnvironmentFile=/etc/pandapages/postgresql-backup.env
ExecStart=/opt/pandapages/scripts/postgresql-backup-healthcheck.sh
TimeoutStartSec=1min
UMask=0077
LogsDirectory=pandapages-postgresql-backup
LogsDirectoryMode=0700
StandardOutput=append:/var/log/pandapages-postgresql-backup/operations.log
StandardError=append:/var/log/pandapages-postgresql-backup/operations.log
StateDirectory=pandapages-postgresql-backup
StateDirectoryMode=0700
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
RestrictAddressFamilies=AF_UNIX
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
